# files/etc/systemd/system/user@.service.d/faster-shutdown.conf
# User session teardown should also cap at 5s.
# Mirrors omarchy default/systemd/user@.service.d/faster-shutdown.conf.
[Service]
TimeoutStopSec=5s
