# files/usr/lib/systemd/system/sddm.service.d/10-wait-for-flatpak-preinstall.conf
# Gate the first graphical login on the Flatpak preinstall service so the
# user's session doesn't come up before Zen/Chromium/Signal/Obsidian/Papers/
# Thunderbird exist — half the Omarchy bindings and mimeapps.list resolve to
# those Flatpak IDs. `Wants=` (not `Requires=`) keeps offline boots moving:
# if network-online times out and flatpak-preinstall fails, sddm still
# starts with broken defaults instead of hanging the boot. Cost on a normal
# online first boot is the install duration (~a minute); subsequent boots
# hit the manifest-hash early exit in the helper and cost ~0.
[Unit]
Wants=flatpak-preinstall.service
After=flatpak-preinstall.service
