# files/etc/systemd/system/plocate-updatedb.service.d/ac-only.conf
# Don't burn battery indexing when unplugged.
# Mirrors omarchy install/config/plocate-ac-only.sh.
[Unit]
ConditionACPower=true
